# files/etc/init.d/99-ath9k-set-txq-memory-limit
#!/bin/sh /etc/rc.common

START=99

start() {
	iw phy phy1 set txq memory_limit 8388608
}
